# app/config/defaults/etc/udev/rules.d/70-bluealsa.rules
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/srv/http/command/bluealsa-config.sh"
